# services/clusterodm.service
[Unit]
Description=Start ClusterODM Service

[Service]
Type=simple
PIDFile=/run/clusterodm.pid
User=odm
Group=odm
WorkingDirectory=/www
ExecStart=/usr/bin/node index.js
ExecStop=/bin/kill -s QUIT $MAINPID
Restart=always

[Install]
WantedBy=multi-user.target
